# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:gateway]
command=gunicorn -b 127.0.0.1:5000 -w 2 --timeout 120 gateway_app:app
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:marketing]
command=bash -c "cd /app && bash /app/ops/run_marketing_migrations.sh && exec gunicorn -b 127.0.0.1:8002 -w 2 --timeout 120 dashboard.app:app"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:producer]
command=bash -c "cd /app && bash /app/ops/run_producer_migrations.sh && cd /app/Producer && python manage.py loadinitialdata && exec gunicorn -b 127.0.0.1:8001 -w 2 --timeout 120 logic_service.wsgi"
environment=DJANGO_SETTINGS_MODULE="logic_service.settings.docker",RUNNING_IN_DOCKER="1",FORCE_SCRIPT_NAME="/producer",SCRIPT_NAME="/producer",STATIC_URL="/producer/static/"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
